# docker/app/supervisord.conf
; nodaemon=true, va a correr en foreground y este a su vez va a ejecutar a los programs listados abajo
[supervisord]
nodaemon=true

; ejecuta el program nginx, manda los logs al logger de docker
[program:nginx]
command=nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ejecuta el program php, manda los logs al logger de docker
[program:php-fpm]
command=php-fpm7.4
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0